# server/systemd/tk-gps-firewall-rollback.timer
[Unit]
Description=TK GPS Firewall automatische Rücksetzung

[Timer]
OnActiveSec=5min
AccuracySec=1s
Unit=tk-gps-firewall-rollback.service

[Install]
WantedBy=timers.target
