# server/systemd/tk-gps-firewall-rollback.service
[Unit]
Description=TK GPS Firewall Rücksetzung

[Service]
Type=oneshot
ExecStart=/var/lib/tk-gps-firewall/rollback.sh
