# server/systemd/tk-nightly-reboot.service
[Unit]
Description=TK nächtliches Backup mit anschließendem Neustart
Wants=network-online.target
After=network-online.target postgresql.service traccar.service
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/tk-nightly-maintenance backup
ExecStartPost=/usr/bin/systemctl reboot
TimeoutStartSec=3h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
